# deploy/systemd/algorithm-vault-worker.service
[Unit]
Description=Algorithm Vault dedicated worker
After=network-online.target postgresql.service algorithm-vault.service
Wants=network-online.target

[Service]
Type=simple
User=tradingbot
Group=tradingbot
WorkingDirectory=/opt/algorithm-vault/current
EnvironmentFile=/etc/algorithm-vault/algorithm-vault.env
Environment=WORKER_MODE=worker
Environment=WORKER_PROCESS_CONFIGURED=true
ExecStartPre=/opt/algorithm-vault/current/scripts/production_bootstrap.sh
ExecStart=/opt/algorithm-vault/current/.venv/bin/python -m app.workers.runner
Restart=on-failure
RestartSec=5
TimeoutStopSec=45
KillSignal=SIGTERM
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/opt/algorithm-vault/current /var/log/algorithm-vault

[Install]
WantedBy=multi-user.target
